# src/utils/recap.service.in
[Unit]
Description=recap
Documentation=man:recap(5) man:recap(8)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=@BINDIR@/recap
