# src/utils/recap-onboot.timer
[Unit]
Description=Timer for recap-onboot

[Timer]
OnBootSec=10s

[Install]
WantedBy=timers.target
