# src/utils/recaplog.timer
[Unit]
Description=Timer for recaplog

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
