# src/utils/recap.service.in
[Unit]
Description=recap
Documentation=man:recap.conf(5) man:recap(8)
After=network-online.target
Wants=network-online.target

[Service]
Environment=BINDIR=/usr/local/bin
Type=oneshot
ExecStart=/bin/sh -c "${BINDIR}/recap"
